# supervisor-eidr-connect.conf
[program:eidr-connect]
command=/bin/bash /eidr-connect.sh
user=meteor
autostart=true
autorestart=unexpected
