# supervisor-eidr-connect.conf
[program:eidr-connect]
command=meteor --allow-superuser --production
user=meteor
autostart=true
autorestart=unexpected
