# supervisor-eidr-connect.conf
[program:eidr-connect]
command=meteor --allow-superuser --production
directory=/eidr-connect/
autostart=true
autorestart=unexpected
